# backend/src/video_processor/supervisord.conf
[supervisord]
nodaemon=true
user=appuser
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:api]
command=python main.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/tmp/api.err.log
stdout_logfile=/tmp/api.out.log
environment=PYTHONUNBUFFERED=1

[program:worker]
command=python worker.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/tmp/worker.err.log
stdout_logfile=/tmp/worker.out.log
environment=PYTHONUNBUFFERED=1

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock 